# skeleton/gunicorn-formhub.conf
description "upstart configuration for gunicorn fh"

start on net-device-up
stop on shutdown

respawn
respawn limit 10 5
setuid "!SERVER_USER!"
setgid "!SERVER_USER!"
exec  !GUNICORN_SHELL_SCRIPT!
